# systemd/dbus-org.nemo.ssu.service
[Unit]
Description=SSU service

[Service]
Type=dbus
BusName=org.nemo.ssu
ExecStart=/usr/bin/ssud

